# templates/general.timer.tmpl
[Unit]
Description=$Description
Requires=$Requiers
[Timer]
Unit=$Unit
OnUnitInactiveSec=$Minutes
AccuracySec=1s
[Install]
WantedBy=$WantedBy
