# templates/general.timer.tmpl
[Unit]
Description=$Description
Requires=$Requires
[Timer]
Unit=$Unit
OnBootSec=5min  # delay after boot
OnUnitInactiveSec=$Timer
AccuracySec=1s
[Install]
WantedBy=timers.target
